# assets/systemd/dynsix.service
[Unit]
Description=IPv6 dyndns client

[Service]
ExecStart=/usr/local/bin/dynsix
